# etc/udev/rules.d/99-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/home/kevin/.local/bin/brightness-control clear"
