# home/dot_config/systemd/user/battery-alert.service.tmpl
[Unit]
Description=Desktop alert warning of low remaining battery

[Service]
Type=oneshot
ExecStart={{ .chezmoi.homeDir }}/.local/bin/check-low-battery

[Install]
WantedBy=graphical.target
